# fvpn/systemd/fvpn-agent.service
[Unit]
Description=FVPN Worker Agent API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/fvpn
EnvironmentFile=/etc/fvpn/fvpn-agent.env
ExecStart=/opt/fvpn/venv/bin/python /opt/fvpn/run_agent.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
